# deploy/supervisor.conf
[program:task_consumer_start]
command=php /var/www/hw21/current/app/bin/console task-consumer:start
process_name=task_consumer_start_%(process_num)02d
numprocs=1
directory=/tmp
autostart=true
autorestart=true
startsecs=3
startretries=10
user=www-data
redirect_stderr=false
stdout_logfile=/var/www/hw21/current/app/var/log/supervisor.task_consumer_start.out.log
stdout_capture_maxbytes=1MB
stderr_logfile=/var/www/hw21/current/app/var/log/supervisor.task_consumer_start.error.log
stderr_capture_maxbytes=1MB
